# argos-sunpos.service
[Unit]
Description=Argos Solar Position Service
After=network.target

[Service]
Type=simple
ExecStart=%h/.local/bin/argos-sunpos-service.py
Restart=on-failure
RestartSec=10
StandardInput=null
StandardOutput=journal
StandardError=journal

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ProtectHostname=true
ProtectKernelTunables=true
ProtectControlGroups=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictRealtime=true
CapabilityBoundingSet=
AmbientCapabilities=
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
SystemCallFilter=~@privileged

[Install]
WantedBy=default.target
